# private_dot_config/autostart/private_fusuma.desktop
[Desktop Entry]
Exec=/home/jason/bin/fusuma.sh -d
Icon=dialog-scripts
Name=Fusuma
Path=
Type=Application
X-KDE-AutostartScript=true
